# prototype/main/im-config/fcitx5/fcitx5-chewing/asset/overlay/etc/profile.d/input-method.sh


##
## ## Arch Wiki
##
## * [Fcitx5](https://wiki.archlinux.org/title/Fcitx5)
## * [Fcitx5](https://wiki.archlinuxcn.org/zh-tw/Fcitx5)
##
##


export GTK_IM_MODULE=fcitx
export QT_IM_MODULE=fcitx
export XMODIFIERS=@im=fcitx
export SDL_IM_MODULE=fcitx
export INPUT_METHOD=fcitx
export GLFW_IM_MODULE=ibus
